# common/.local/etc/profile.d/desktop.sh
# This file should be overwritten by dotfile packages to set environment
# variables specific to a desktop environment.  That way, only one such
# file will exist, even if multiple desktop environment dotfile packages
# are installed simultaneously, allowing the preferred configuration to
# be easily specified at installation.
